# scripts/app.conf
[program:nacciolyapp]
command=node /appdir/server.js
directory=/appdir 
user=ec2-user
autorestart=true
stderr_logfile=/var/log/meuapp.err.log
stdout_logfile=/var/log/meuapp.out.log